# supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=python manage.py runserver 0.0.0.0:8000

[program:cron]
command=cron -f